# container-assets/systemd/nix-daemon.socket.d/override.conf
[Unit]
Requires=discobot-nix-restore.service
After=discobot-nix-restore.service
