# vm-assets/systemd/systemd-timesyncd.service.d/role.conf
[Unit]
# systemd-timesyncd persists state under /var/lib/systemd/timesync, so it
# must wait for init-var.service to mount the writable /var disk first.
Requires=init-var.service
After=init-var.service

[Service]
EnvironmentFile=-/etc/default/discobot-wsl
ExecCondition=/usr/local/bin/check-wsl-role.sh platform non-wsl
