# vm-assets/systemd/init-var.service
[Unit]
Description=Prepare persistent /var storage
DefaultDependencies=no
After=local-fs.target systemd-journald.socket
# Run in sysinit before services that need writable state under /var/lib.
Before=sysinit.target docker.service containerd.service
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/default/discboeing-wsl
ExecStart=/usr/local/bin/init-var.sh
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=sysinit.target
