# vm-assets/systemd/init-var.service
[Unit]
Description=Prepare persistent /var storage
After=local-fs.target
# Run before Docker services (don't include journald to avoid cycles)
Before=docker.service containerd.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/default/discobot-wsl
ExecStart=/usr/local/bin/init-var.sh
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
