# container-assets/systemd/discobot-vscode-backend.service
[Unit]
Description=Discobot VS Code backend (code-server)
After=discobot-setup.service
Requires=discobot-setup.service

[Service]
Type=simple
User=discobot
Group=discobot
WorkingDirectory=/home/discobot/workspace
Environment=HOME=/home/discobot
EnvironmentFile=-/run/discobot/agent-env
ExecStartPre=/usr/bin/mkdir -p /home/discobot/.local/share/discobot-code-server/User /home/discobot/.local/share/discobot-code-server/extensions
ExecStartPre=/bin/bash -lc 'cp -r --update=none /opt/discobot/code-server-defaults/User/. /home/discobot/.local/share/discobot-code-server/User/'
ExecStartPre=/bin/bash -lc 'cp -r --update=none /opt/discobot/code-server-defaults/extensions/. /home/discobot/.local/share/discobot-code-server/extensions/'
ExecStartPre=/bin/bash -lc 'rm -rf /home/discobot/.local/share/discobot-code-server/extensions/discobot.discobot-theme && cp -r /opt/discobot/code-server-defaults/extensions/discobot.discobot-theme /home/discobot/.local/share/discobot-code-server/extensions/'
ExecStart=/usr/bin/code-server --bind-addr 127.0.0.1:13338 --auth none --disable-telemetry --disable-update-check --user-data-dir /home/discobot/.local/share/discobot-code-server --extensions-dir /home/discobot/.local/share/discobot-code-server/extensions /home/discobot/workspace
Restart=on-failure
RestartSec=2

StandardOutput=journal+console
StandardError=journal+console
