# container-assets/systemd/websockify.service
[Unit]
Description=WebSocket to VNC proxy with desktop resize support
Requires=xvfb.service
After=xvfb.service

[Service]
Type=simple
User=discboeing
Group=discboeing
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
CapabilityBoundingSet=
ExecStart=/usr/local/bin/discboeing-vnc-websockify
Restart=on-failure
RestartSec=2
