# container-assets/systemd/x11vnc.socket
[Unit]
Description=x11vnc VNC server socket activation on localhost port 5900

[Socket]
# x11vnc is unauthenticated; keep the raw VNC socket container-local and let
# the local websockify wrapper mediate browser access.
ListenStream=127.0.0.1:5900
Accept=true

[Install]
WantedBy=sockets.target
