# container-assets/systemd/xvfb.service
[Unit]
Description=Xvfb virtual X11 framebuffer on display :0
# Only started on-demand via x11-display.socket or x11vnc.socket dependencies
# Xvfb binds to TCP port 6000 (display :0) so systemd-socket-proxyd can forward
# X11 unix socket connections to it
Wants=openbox.service

[Service]
Type=simple
# Remove stale X11 lock file and socket that would prevent Xvfb from starting
ExecStartPre=-/bin/rm -f /tmp/.X11-unix/X0 /tmp/.X0-lock
ExecStart=/usr/bin/Xvfb :0 -screen 0 1280x1024x24 -ac -listen tcp -nolisten local
Restart=on-failure
RestartSec=2
