# container-assets/systemd/nix-daemon.service.d/override.conf
[Unit]
Requires=discboeing-nix-restore.service
After=discboeing-nix-restore.service

[Service]
ExecStartPre=/usr/local/bin/restore-nix-staging
