# vm-assets/systemd/systemd-networkd-wait-online.service.d/role.conf
[Service]
EnvironmentFile=-/etc/default/discboeing-wsl
ExecCondition=/usr/local/bin/check-wsl-role.sh platform non-wsl
